# machine-builders/ubuntu-mirror/fs/etc/cron.daily/apt-mirror
#!/bin/sh

# Disabled for now because this corrupts the apt db
# apt-mirror
